# archivos_SH/serviceNode.service
[Unit]
Description=Mi Servicio Node.js

[Service]
ExecStart=/usr/local/bin/ejecutarNode.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
